# infrastructure/ansible/net-status/roles/net-status/files/netstatus.conf
respawn
 
start on runlevel [2345]
stop on runlevel [06]
 
script
su - netstatus -s /bin/bash -c "/home/netstatus/gen_status_loop.py 2>&1" >>  /home/netstatus/netstatus.log
end script
